# threat_intel_report.json
{
  "IoCs": {
    "Domains": [
      "www.recordedfuture.com",
      "protonmail.com",
      "recordedfuture.com"
    ],
    "IP addresses": []
  },
  "TTPs": {
    "Tactics": [
      [
        "TA0002",
        "Execution"
      ],
      [
        "TA0040",
        "Impact"
      ],
      [
        "TA0007",
        "Discovery"
      ],
      [
        "TA0009",
        "Collection"
      ]
    ],
    "Techniques": [
      [
        "T1053.005",
        "Scheduled Task"
      ],
      [
        "T1116",
        "Code Signing"
      ],
      [
        "T1587.001",
        "Malware"
      ],
      [
        "T1553.002",
        "Code Signing"
      ],
      [
        "T1153",
        "Source"
      ],
      [
        "T1505.003",
        "Web Shell"
      ],
      [
        "T1583.004",
        "Server"
      ],
      [
        "T1589.002",
        "Email Addresses"
      ],
      [
        "T1588.001",
        "Malware"
      ],
      [
        "T1059.001",
        "PowerShell"
      ],
      [
        "T1588.002",
        "Tool"
      ],
      [
        "T1592.002",
        "Software"
      ],
      [
        "T1589.001",
        "Credentials"
      ],
      [
        "T1100",
        "Web Shell"
      ],
      [
        "T1485",
        "Data Destruction"
      ],
      [
        "T1584.004",
        "Server"
      ],
      [
        "T1053.002",
        "At"
      ],
      [
        "T1086",
        "PowerShell"
      ],
      [
        "T1001.001",
        "Junk Data"
      ]
    ]
  },
  "Threat Actor(s)": [
    "Sandworm"
  ],
  "Malware": [
    {
      "Name": "Unknown",
      "tags": [
        "Detected by 26/92 engines",
        "Detected by 30/92 engines",
        "Detected by 31/92 engines",
        "Detected by 59/92 engines",
        "Detected by 26/92 engines"
      ],
      "sha256": "b01e0c6ac0b8bcde145ab7b68cf246deea9402fa7ea3aede7105f7051fe240c1",
      "sha1": "1ae7556dfacd47d9efbe79be974661a5a6d6d923"
    },
    {
      "Name": "H",
      "tags": [
        "Detected by 59/92 engines",
        "Detected by 57/92 engines",
        "Detected by 58/92 engines",
        "Detected by 57/92 engines",
        "Detected by 60/92 engines"
      ],
      "sha256": "1bc44eef75779e3ca1eefb8ff5a64807dbc942b1e4a2672d77b9f6928d292591"
    }
  ],
  "Targeted Entities": [
    "finance and government contractor organizations in Ukraine",
    "Ukrainian websites earlier that day",
    "a Microsoft Exchange server on December",
    "Ukraine",
    "Ukrainian organizations"
  ]
}